# verbatimd.plist.in
<?xml version="1.0" encoding="utf-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
   "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>local.verbatimd</string>

		<key>ProgramArguments</key>
		<array>
			<string>~/.local/bin/verbatimd</string>
			<string>--host</string>
			<string>127.0.0.1</string>
			<string>--port</string>
			<string>5959</string>
			<string>--log-level</string>
			<string>info</string>
		</array>

		<key>RunAtLoad</key>
		<true />

		<key>KeepAlive</key>
		<true />

		<key>StandardOutPath</key>
		<string>~/Library/Logs/verbatimd/stdout.log</string>

		<key>StandardErrorPath</key>
		<string>~/Library/Logs/verbatimd/stderr.log</string>
	</dict>
</plist>
